# gluon/gluon-setup-mode/files/lib/gluon/setup-mode/rc.d/S90geolocate-ffgt
#!/bin/sh /etc/rc.common

if [ -e /lib/gluon/gluon-luci-geolocate/geolocate.sh ]; then
    ((/lib/gluon/gluon-luci-geolocate/geolocate.sh)&)
fi
